# files/etc/uci-defaults/99-daed-test-network.template
#!/bin/sh

# === Network (LAN) ===
uci set network.lan.ipaddr='{{ LAN_IP }}'
uci set network.lan.netmask='255.255.255.0'
uci set network.lan.delegate='0'

# === SSH (dropbear) - 全网口访问 ===
uci set dropbear.main.enable='1'
uci set dropbear.main.PasswordAuth='on'
uci set dropbear.main.RootPasswordAuth='on'
uci set dropbear.main.Port='22'
uci delete dropbear.main.Interface 2>/dev/null || true

# === Root password ===
echo -e "{{ ROOT_PASSWORD }}\n{{ ROOT_PASSWORD }}" | passwd root

# === IPv6 ===
uci set network.lan.ipv6='1'
uci set network.wan.ipv6='1'
uci set dhcp.lan.ra='server'
uci set dhcp.lan.dhcpv6='server'
uci set dhcp.lan.radvd='1'

# === PPPoE 拨号 ===
if [ "{{ ENABLE_PPPOE }}" = "yes" ]; then
  uci set network.wan.proto='pppoe'
  uci set network.wan.username='{{ PPPOE_ACCOUNT }}'
  uci set network.wan.password='{{ PPPOE_PASSWORD }}'
  uci set network.wan.ipv6='1'
fi

# === WAN 防火墙入站规则 ===
if [ "{{ ENABLE_FW_WAN }}" = "yes" ]; then
  uci add firewall rule
  uci set firewall.@rule[-1].name='Allow-WAN-Inbound'
  uci set firewall.@rule[-1].src='wan'
  uci set firewall.@rule[-1].target='ACCEPT'
fi

# === Docker 防火墙规则 ===
if [ "{{ ENABLE_DOCKER }}" = "yes" ]; then
  uci add firewall zone
  uci set firewall.@zone[-1].name='docker'
  uci set firewall.@zone[-1].input='ACCEPT'
  uci set firewall.@zone[-1].output='ACCEPT'
  uci set firewall.@zone[-1].forward='REJECT'

  uci add firewall forwarding
  uci set firewall.@forwarding[-1].src='docker'
  uci set firewall.@forwarding[-1].dest='wan'

  uci add firewall forwarding
  uci set firewall.@forwarding[-1].src='lan'
  uci set firewall.@forwarding[-1].dest='docker'
fi

# === Android TV DNS 修复 ===
uci add dhcp domain
uci set dhcp.@domain[-1].name='time.android.com'
uci set dhcp.@domain[-1].ip='216.239.36.1'

uci commit network
uci commit firewall
uci commit dhcp

# === Configure uhttpd (internal, for LuCI CGI) ===
/etc/init.d/nginx stop 2>/dev/null || true
/etc/init.d/uhttpd stop 2>/dev/null || true

uci set uhttpd.main.listen_http='127.0.0.1:8080'
uci set uhttpd.main.listen_https=''
uci commit uhttpd

# === Self-signed SSL certificate ===
SSL_DIR="/etc/nginx/ssl"
mkdir -p "$SSL_DIR"
if [ ! -f "$SSL_DIR/server.crt" ]; then
  openssl req -x509 -nodes -days 3650 -newkey rsa:2048 \
    -keyout "$SSL_DIR/server.key" \
    -out "$SSL_DIR/server.crt" \
    -subj "/CN={{ SSL_CN }}/O=ImmortalWrt/C=CN" 2>/dev/null
fi

# === nginx configuration ===
cat > /etc/nginx/nginx.conf <<'NGINX'
worker_processes auto;
pid /run/nginx.pid;

events {
    worker_connections 1024;
}

http {
    include /etc/nginx/mime.types;
    default_type application/octet-stream;
    sendfile on;
    keepalive_timeout 65;

    server {
        listen 80 default_server;
        listen [::]:80 default_server;
        server_name _;
        return 301 https://$host$request_uri;
    }

    server {
        listen 443 ssl default_server;
        listen [::]:443 ssl default_server;
        server_name _;

        ssl_certificate /etc/nginx/ssl/server.crt;
        ssl_certificate_key /etc/nginx/ssl/server.key;
        ssl_protocols TLSv1.2 TLSv1.3;
        ssl_ciphers HIGH:!aNULL:!MD5;

        location /op/ {
            alias /www/;
            try_files $uri $uri/ =404;
        }

        location / {
            proxy_pass http://127.0.0.1:8080;
            proxy_redirect http:// $scheme://;
            proxy_set_header Host $host;
            proxy_set_header X-Real-IP $remote_addr;
            proxy_set_header X-Forwarded-For $proxy_add_x_forwarded_for;
            proxy_set_header X-Forwarded-Proto $scheme;
        }
    }
}
NGINX

/etc/init.d/uhttpd enable
/etc/init.d/uhttpd start
/etc/init.d/nginx enable
/etc/init.d/nginx start

/etc/init.d/network restart
/etc/init.d/dropbear enable
/etc/init.d/dropbear restart

exit 0
